# contrib/lighthouse-plugin/build/plugin-server.service
[Unit]
Description=Lighthouse plugin server
Before=kubelet.service

[Service]
Type=notify
EnvironmentFile=-/etc/plugin-server/config
ExecStart=/usr/bin/plugin-server $ARGS
Restart=always

LimitNOFILE=infinity
LimitNPROC=infinity
LimitCORE=infinity
# Uncomment TasksMax if your systemd version supports it.
# Only systemd 226 and above support this version.
TasksMax=infinity

[Install]
WantedBy=multi-user.target
